# outputs/dev_results/v3/DOC003.json
{
  "report_metadata": {
    "report_type": "incident",
    "report_title": "CISA Shares Lessons Learned from an Incident Response Engagement"
  },
  "threat_actor": null,
  "campaign_name": null,
  "target_sectors": [
    "Federal Civilian Executive Branch (FCEB) agency"
  ],
  "indicators_of_compromise": {
    "ip_addresses": [],
    "domains": [],
    "file_hashes": [],
    "urls": [],
    "http_paths": [],
    "tools": [
      {
        "name": "Burp Suite Burp Scanner",
        "evidence_snippet": "CISA detected this scanning activity by analyzing web logs and identifying signatures associated with the tool."
      },
      {
        "name": "RingQ",
        "evidence_snippet": "The cyber threat actors staged RingQ, a publicly available defense evasion tool, for potential use."
      },
      {
        "name": "Stowaway",
        "evidence_snippet": "The cyber threat actors used Stowaway, a publicly available multi-level proxy tool, to establish C2."
      },
      {
        "name": "fscan",
        "evidence_snippet": "CISA identified the use of the fscan tool by analyzing evidence of its output found on disk."
      }
    ],
    "cves": [
      {
        "value": "CVE-2024-36401",
        "evidence_snippet": "The cyber threat actors exploited CVE-2024-36401 for initial access on two GeoServers.",
        "source_context": null
      }
    ],
    "detection_signatures": {
      "yara_rules": [],
      "sigma_rules": []
    }
  },
  "ttps": [
    {
      "mitre_id": "T1595.002",
      "tactic": "Reconnaissance",
      "technique": "Active Scanning: Vulnerability Scanning",
      "evidence_snippet": "The cyber threat actors performed active scanning to identify vulnerabilities they could use for initial access.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1583.003",
      "tactic": "Resource Development",
      "technique": "Acquire Infrastructure: Virtual Private Server",
      "evidence_snippet": "The cyber threat actors gained remote access to the victim’s network using a desktop behind a virtual private server (VPS).",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1190",
      "tactic": "Initial Access",
      "technique": "Exploit Public-Facing Application",
      "evidence_snippet": "The cyber threat actors exploited CVE 2024-36401 on two of the organization’s public-facing GeoServers.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1059.001",
      "tactic": "Execution",
      "technique": "Command and Scripting Interpreter: PowerShell",
      "evidence_snippet": "The cyber threat actors used PowerShell to download a payload.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1202",
      "tactic": "Defense Evasion",
      "technique": "Indirect Command Execution",
      "evidence_snippet": "The cyber threat actors employed indirect command execution via web shells.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1197",
      "tactic": "Persistence",
      "technique": "BITS Jobs",
      "evidence_snippet": "The cyber threat actors abused BITS jobs.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1053.003",
      "tactic": "Persistence",
      "technique": "Scheduled Task/Job: Cron",
      "evidence_snippet": "The cyber threat actors established persistence through cron jobs.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1505.003",
      "tactic": "Persistence",
      "technique": "Server Software Component: Web Shell",
      "evidence_snippet": "The cyber threat actors uploaded web shells for persistence.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1078",
      "tactic": "Persistence",
      "technique": "Valid Accounts",
      "evidence_snippet": "The cyber threat actors used valid accounts for persistence.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1068",
      "tactic": "Privilege Escalation",
      "technique": "Exploitation for Privilege Escalation",
      "evidence_snippet": "The cyber threat actors attempted to exploit CVE-2016-5195 to escalate privileges.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1110",
      "tactic": "Credential Access",
      "technique": "Brute Force",
      "evidence_snippet": "The cyber threat actors used brute force techniques to obtain login credentials for web services.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1087.001",
      "tactic": "Discovery",
      "technique": "Account Discovery: Local Account",
      "evidence_snippet": "The cyber threat actors used cat /etc/passwd to discover local users.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1083",
      "tactic": "Discovery",
      "technique": "File and Directory Discovery",
      "evidence_snippet": "The cyber threat actors used dir c:\\, dir d:\\, dir e:\\, and type c:\\ commands to identify files and directories on the SQL server.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1046",
      "tactic": "Discovery",
      "technique": "Network Service Discovery",
      "evidence_snippet": "The cyber threat actors used fscan to identify SSH listeners and FTP servers.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1057",
      "tactic": "Discovery",
      "technique": "Process Discovery",
      "evidence_snippet": "The cyber threat actors used tasklist on the SQL server.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1018",
      "tactic": "Discovery",
      "technique": "Remote System Discovery",
      "evidence_snippet": "The cyber threat actors performed ping sweeps of hosts within specific subnets.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1082",
      "tactic": "Discovery",
      "technique": "System Information Discovery",
      "evidence_snippet": "The cyber threat actors used cat /etc/redhat-release and cat /etc/os-release commands to get Red Hat Enterprise Linux (RHEL) and Linux operating system information.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1016",
      "tactic": "Discovery",
      "technique": "System Network Configuration Discovery",
      "evidence_snippet": "The cyber threat actors used ipconfig to check GeoServer 1’s and the SQL server’s network configurations.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1049",
      "tactic": "Discovery",
      "technique": "System Network Connections Discovery",
      "evidence_snippet": "The cyber threat actors executed commands such as netstat to obtain a listing of network connections to or from the systems they compromised.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1033",
      "tactic": "Discovery",
      "technique": "System Owner/User Discovery",
      "evidence_snippet": "The cyber threat actors used whoami on the SQL server.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1105",
      "tactic": "Command and Control",
      "technique": "Ingress Tool Transfer",
      "evidence_snippet": "The cyber threat actors used PowerShell and bitsadmin getfile to download payloads.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    },
    {
      "mitre_id": "T1090",
      "tactic": "Command and Control",
      "technique": "Proxy",
      "evidence_snippet": "The cyber threat actors used a connection proxy to direct traffic from their C2 server.",
      "mapping_type": "explicit",
      "mitre_confidence": null,
      "evidence_location": null
    }
  ],
  "summary": "CISA shares lessons learned from an incident response engagement at a U.S. federal civilian executive branch (FCEB) agency, highlighting the importance of prompt patching and preparing for incidents by practicing incident response plans and implementing logging and aggregating logs in a centralized out-of-band location.",
  "timeline": []
}